# ops/systemd/aialra-babbledeck-backup.service
[Unit]
Description=AIALRA BabbleDeck production backup
After=docker.service
Wants=docker.service

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/srv/aialra/apps/codexapp/state/browser-workspaces/2026-07-04-babbledeck
EnvironmentFile=/srv/aialra/config/secrets/babbledeck.env
ExecStart=/usr/bin/bash /srv/aialra/apps/codexapp/state/browser-workspaces/2026-07-04-babbledeck/scripts/backup-production.sh
StandardOutput=append:/srv/aialra/logs/babbledeck/backup.log
StandardError=append:/srv/aialra/logs/babbledeck/backup.log

[Install]
WantedBy=multi-user.target
